# root/etc/services.d/netdata/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/netdata -c /config/netdata.conf -D
